# testing/logrotate/logrotate.cron.daily
#!/bin/sh

/sbin/logrotate /etc/logrotate.conf
