# core/shadow/shadow.cron.daily
#!/bin/sh

# Verify integrity of password and group files
/sbin/pwck -r
/sbin/grpck -r

